# scrn_rotd.service
[Unit]
Description=Auto Screen Rotation Daemon

[Service]
ExecStart=/usr/bin/scrn_rotd $h/.scrn_rotd.conf

[Install]
WantedBy=default.target